# systemd/dawos-agent.service
[Unit]
Description=dawos-agent — PPP router management daemon
Documentation=https://github.com/Cepat-Kilat-Teknologi/dawos-agent
After=network-online.target accel-ppp.service
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=dawos
Group=dawos
EnvironmentFile=-/etc/dawos-agent/agent.env
ExecStart=/opt/dawos-agent/venv/bin/dawos-agent
Restart=always
RestartSec=3
WatchdogSec=30

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=-/etc/accel-ppp.conf -/etc/accel-ppp.d -/etc/accel-nat-egress.nft -/etc/sysctl.d -/etc/nftables.conf -/etc/resolv.conf -/etc/systemd/resolved.conf.d -/etc/dnsmasq.d -/etc/dnsmasq.conf
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dawos-agent

[Install]
WantedBy=multi-user.target
