# systemd/dawos-agent.service
[Unit]
Description=dawos-agent — PPP router management daemon
Documentation=https://github.com/Cepat-Kilat-Teknologi/dawos-agent
After=network-online.target accel-ppp.service
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=dawos
Group=dawos
EnvironmentFile=-/etc/dawos-agent/agent.env
ExecStart=/opt/dawos-agent/venv/bin/dawos-agent
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=-/etc/accel-ppp.conf -/etc/accel-ppp.d -/etc/accel-nat-egress.nft -/etc/sysctl.d -/etc/nftables.conf -/etc/resolv.conf
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dawos-agent

[Install]
WantedBy=multi-user.target
